# QoL/.config/systemd/user/battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/mark/.local/bin/battery-notify.sh
