# QoL/.config/systemd/user/battery.timer
[Unit]
Description=Timer to check battery status

[Timer]
OnActiveSec=20min
OnUnitActiveSec=20min

[Install]
WantedBy=graphical.target
